# etc/systemd/system/before-suspend.service
[Unit]
Description=Before suspend script
Before=sleep.target

[Service]
Type=simple
ExecStart=/bin/sh -c "light -O"

[Install]
WantedBy=sleep.target
